# alpine-php-mysql/0.3/supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/LogFiles/supervisor/supervisord.log

[program:sshd]
command=rc-service sshd start

[program:apache2]
command=/usr/sbin/httpd -D FOREGROUND